# deploy/sentinel-agent.service
# systemd unit for the Sentinel endpoint agent (production install).
#   sudo cp deploy/sentinel-agent.service /etc/systemd/system/
#   sudo install -m0755 bin/sentinel-agent /usr/local/bin/sentinel-agent
#   sudo systemctl enable --now sentinel-agent
[Unit]
Description=Sentinel Endpoint Agent (EDR/DLP)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
# Enrollment + server config (override in /etc/sentinel/agent.env)
# Expected env includes SENTINEL_SERVER, SENTINEL_ENROLL_TOKEN, and, when mTLS is enabled,
# SENTINEL_AGENT_TLS_CA, SENTINEL_AGENT_TLS_CERT, SENTINEL_AGENT_TLS_KEY.
EnvironmentFile=-/etc/sentinel/agent.env
ExecStart=/usr/local/bin/sentinel-agent
Restart=always
RestartSec=5
# Respawn even on a "clean" exit code or a SIGKILL/SIGTERM that didn't come from a manual
# `systemctl stop` — defeats `systemctl kill`, `kill -9`, self-exit, and crashes. (A manual
# `systemctl stop`/`disable` is still honored by systemd; that path is covered by the kernel
# kill-deny + a loud tamper alert, not by fighting the init system, which would also break the
# verified self-update restart and in-place upgrades.)
RestartForceExitStatus=0 1 2 SIGTERM SIGKILL
# Always-on, survives reboots; restarts on crash.
StateDirectory=sentinel
Environment=SENTINEL_STATE=/var/lib/sentinel/agent.json

# Capabilities the agent genuinely needs: nftables isolation/upload-block (NET_ADMIN),
# eBPF/perf capture (BPF/PERFMON), DNS sniffing (NET_RAW), USB module control (SYS_MODULE),
# process containment/inspection (KILL/SYS_PTRACE), and root-owned file reads/quarantine.
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_MODULE CAP_BPF CAP_PERFMON CAP_NET_RAW CAP_KILL CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_MODULE CAP_BPF CAP_PERFMON CAP_NET_RAW CAP_KILL CAP_SYS_PTRACE CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE

# Self-protection + sandboxing. Keep /usr/local/bin writable so verified self-update can
# atomically stage and replace /usr/local/bin/sentinel-agent; the rest of /usr remains read-only.
NoNewPrivileges=true
OOMScoreAdjust=-1000
ProtectSystem=full
ReadWritePaths=/var/lib/sentinel /usr/local/bin
ProtectHome=read-only
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
PrivateTmp=true
SystemCallArchitectures=native

# Intentionally NOT set — each would break a core EDR capability:
#   ProtectKernelModules    → agent unloads USB storage modules and may load eBPF/kernel helpers
#   MemoryDenyWriteExecute  → eBPF JIT can need W^X exemptions on some hosts
#   PrivateDevices          → USB / removable-media monitoring needs /dev
#   RestrictAddressFamilies → AF_NETLINK/packet sockets are required for process/net telemetry
#   ProtectControlGroups    → cgroup-freeze response needs cgroup access
#   ProtectKernelTunables   → posture collection reads kernel/sysctl state
#   SystemCallFilter        → too broad to constrain safely for a full-system monitor

[Install]
WantedBy=multi-user.target
